# deploy/systemd/dream-diary-bot.service
[Unit]
Description=Dream Diary Telegram Bot
After=network-online.target mongod.service
Wants=network-online.target

[Service]
Type=simple
User=__BOT_USER__
WorkingDirectory=__BOT_DIR__
EnvironmentFile=-%h/.config/dreams.env
EnvironmentFile=-__BOT_DIR__/.env
ExecStart=__BOT_DIR__/.venv/bin/python __BOT_DIR__/main.py
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=20
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
